# meta-user/meta-user/meta-user/recipes-apps/init-script/files/init-script.service
[Unit]
Description=init-script

[Service]
ExecStart=/usr/bin/init-script
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target

